# netflix/provisioning/conf/app/supervisor/app.conf
[program:app]
command=/opt/groovy/bin/groovy /opt/app.groovy
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log
