# etcd/provisioning/conf/etcd/supervisor.conf
[program:etcd]
command=/opt/etcd/etcd
autostart=true
autorestart=true
stderr_logfile=/var/log/etcd.err.log
stdout_logfile=/var/log/etcd.out.log
