# provisioning/conf/db/supervisor/nerve.conf
[program:nerve]
command=/usr/local/bin/nerve -c /etc/nerve.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/nerve.err.log
stdout_logfile=/var/log/nerve.out.log
