# provisioning/conf/app/supervisor/synapse.conf
[program:synapse]
command=/usr/local/bin/synapse -c /etc/synapse.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/synapse.err.log
stdout_logfile=/var/log/synapse.out.log
